# .config/systemd/user/rcloneUtil.service
[Unit]
Description=RClone sync and mount utility

[Service]
Type=simple
ExecStart=/home/alexycn/bin/rcloneUtil.sh

[Install]
WantedBy=default.target
